# supervisor.unison.conf
[program:unison]
command=/unison.sh
user=root
redirect_stderr=true
