# supervisor.fswatch.conf
[program:fswatch]
command=/fswatch.sh
user=root
redirect_stderr=true
